# deploy/backend.service
[Unit]
Description=Power Trading Platform - Backend (FastAPI)
After=network.target

[Service]
User=kartik-gcp
Group=kartik-gcp
WorkingDirectory=/home/kartik-gcp/makeathon-20
Environment="PATH=/home/kartik-gcp/makeathon-20/.venv/bin:/usr/local/bin:/usr/bin"
Environment="SCHEDULER_ENABLED=true"
Environment="SCHEDULER_TIMEZONE=Asia/Kolkata"
Environment="SCHEDULER_CRON_HOUR=1"
Environment="SCHEDULER_CRON_MINUTE=10"
Environment="SCHEDULER_SCRAPE_LOOKBACK_DAYS=1"
Environment="SCHEDULER_RETRAIN_ENABLED=true"
Environment="SCHEDULER_TRAIN_SEGMENTS=DAM,RTM,TAM"
Environment="SCHEDULER_TAM_REGION=NR"
Environment="SCHEDULER_RUN_ON_STARTUP=false"
ExecStart=/home/kartik-gcp/makeathon-20/.venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
