# deploy/backend.service
[Unit]
Description=Power Trading Platform - Backend (FastAPI)
After=network.target

[Service]
User=kartik-gcp
Group=kartik-gcp
WorkingDirectory=/home/kartik-gcp/makeathon-20
Environment="PATH=/home/kartik-gcp/makeathon-20/.venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/home/kartik-gcp/makeathon-20/.venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
